# {{cookiecutter.repo_name}}/build_templates/supervisord.conf.in
[unix_http_server]
file = ${buildout:directory}/var/supervisord.sock

; for unix/linux you can use files instead of TCP ports
; comment [inet_http_server] and uncomment [unix_http_server] section
; and change [supervisorctl] section to use file socket than tcp
;[unix_http_server]
;file = ${buildout:directory}/var/supervisord.sock

[supervisord]
logfile = ${buildout:directory}/var/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid
nodaemon = false

[supervisorctl]
serverurl = unix://${buildout:directory}/var/supervisord.sock
prompt = django

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[include]
files = ${buildout:directory}/etc/*-supervisor.conf

[program:{{ cookiecutter.project_name }}]
directory = ${buildout:directory}
command = ${buildout:directory}/bin/uwsgi --xml ${buildout:directory}/parts/uwsgi/uwsgi.xml
autostart = true
autorestart = true
stopsignal=INT
